# scripts-for-other-scripts-only/drawpile-srv.socket
[Unit]
Description=Drawpile dedicated server socket
After=network.target

[Socket]
ListenStream=27750
# Adding a second port activates the web admin API:
ListenStream=27780
# Adding a third port activates the web client:
ListenStream=27751

[Install]
WantedBy=sockets.target

